# data/nm-sudo.service.in
[Unit]
Description=Network Manager Sudo Helper

[Service]
Type=dbus
BusName=org.freedesktop.nm.sudo
ExecStart=@libexecdir@/nm-sudo

[Install]
Alias=dbus-org.freedesktop.nm.sudo.service

[Service]
AmbientCapabilities=
CapabilityBoundingSet=~CAP_(CHOWN|FSETID|SETFCAP)
CapabilityBoundingSet=~CAP_(DAC_*|FOWNER|IPC_OWNER)
CapabilityBoundingSet=~CAP_AUDIT_*
CapabilityBoundingSet=~CAP_BLOCK_SUSPEND
CapabilityBoundingSet=~CAP_IPC_LOCK
CapabilityBoundingSet=~CAP_KILL
CapabilityBoundingSet=~CAP_LEASE
CapabilityBoundingSet=~CAP_LINUX_IMMUTABLE
CapabilityBoundingSet=~CAP_MAC_*
CapabilityBoundingSet=~CAP_MKNOD
CapabilityBoundingSet=~CAP_NET_(BIND_SERVICE|BROADCAST|RAW)
CapabilityBoundingSet=~CAP_NET_ADMIN
CapabilityBoundingSet=~CAP_SET(UID|GID|PCAP)
CapabilityBoundingSet=~CAP_SYSLOG
CapabilityBoundingSet=~CAP_SYS_(NICE|RESOURCE)
CapabilityBoundingSet=~CAP_SYS_ADMIN
CapabilityBoundingSet=~CAP_SYS_BOOT
CapabilityBoundingSet=~CAP_SYS_CHROOT
CapabilityBoundingSet=~CAP_SYS_MODULE
CapabilityBoundingSet=~CAP_SYS_PACCT
CapabilityBoundingSet=~CAP_SYS_PTRACE
CapabilityBoundingSet=~CAP_SYS_RAWIO
CapabilityBoundingSet=~CAP_SYS_TIME
CapabilityBoundingSet=~CAP_SYS_TTY_CONFIG
CapabilityBoundingSet=~CAP_WAKE_ALARM
Delegate=
DeviceAllow=
IPAddressDeny=
KeyringMode=
LockPersonality=
MemoryDenyWriteExecute=
NoNewPrivileges=
NotifyAccess=
PrivateDevices=
PrivateMounts=
PrivateNetwork=
PrivateTmp=
PrivateUsers=
ProcSubset=
ProtectClock=
ProtectControlGroups=
ProtectHome=
ProtectHostname=
ProtectKernelLogs=
ProtectKernelModules=
ProtectKernelTunables=
ProtectProc=
ProtectSystem=
RemoveIPC=
RestrictAddressFamilies=~AF_(INET|INET6)
RestrictAddressFamilies=~AF_NETLINK
RestrictAddressFamilies=~AF_PACKET
RestrictAddressFamilies=~AF_UNIX
RestrictAddressFamilies=~…
RestrictNamespaces=~CLONE_NEWCGROUP
RestrictNamespaces=~CLONE_NEWIPC
RestrictNamespaces=~CLONE_NEWNET
RestrictNamespaces=~CLONE_NEWNS
RestrictNamespaces=~CLONE_NEWPID
RestrictNamespaces=~CLONE_NEWUSER
RestrictNamespaces=~CLONE_NEWUTS
RestrictRealtime=
RestrictSUIDSGID=
RootDirectory=/RootImage=
SupplementaryGroups=
SystemCallArchitectures=
SystemCallFilter=~@clock
SystemCallFilter=~@cpu-emulation
SystemCallFilter=~@debug
SystemCallFilter=~@module
SystemCallFilter=~@mount
SystemCallFilter=~@obsolete
SystemCallFilter=~@privileged
SystemCallFilter=~@raw-io
SystemCallFilter=~@reboot
SystemCallFilter=~@resources
SystemCallFilter=~@swap
UMask=
User=/DynamicUser=
